# com.rutinerad.recept.plist
  <?xml version="1.0" encoding="UTF-8"?>
  <!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
  <plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.rutinerad.recept</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/samuel/Server/mela-browser/.venv/bin/python</string>
      <string>-m</string>
      <string>mela_browser</string>
      <string>--port</string>
      <string>15482</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/samuel/Server/mela-browser</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/mela-browser.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mela-browser.err</string>
  </dict>
  </plist>